# data/conf/examples/drhood.conf
description "rhood daemon"

start on (filesystem)
stop on runlevel [016]

respawn
console log
setuid www-data
setgid www-data

chdir /home/ubuntu/rhood/rhood-www

script
    /home/ubuntu/rhood/cmd/rhood//rhood \
    --bind-proxy="localhost:8081" \
    --bind-web="10.0.3.155:8090" \
    --public-address="10.0.3.155:90" \
    --public-address-ws="10.0.3.155:8090"

end script
