# data/conf/examples/drhood.conf
description "rhood daemon"

start on (filesystem)
stop on runlevel [016]

respawn
console log
setuid www-data
setgid www-data

chdir /home/ubuntu/rhood/rhood-www

script
    /home/ubuntu/rhood/cmd/rhood//rhood
end script